# deploy/quadlets/lappbot-postgres.container
[Unit]
Description=Lappbot Postgres Database
After=network-online.target

[Container]
Image=postgres:15-alpine
# Map specific required usage variables if they differ from standard .env
# Note: Ensure .env contains POSTGRES_USER, POSTGRES_PASSWORD, POSTGRES_DB
EnvironmentFile=%h/lappbot/.env
Volume=lappbot-postgres-data:/var/lib/postgresql/data
Network=lappbot.network
ContainerName=lappbot-postgres

[Service]
Restart=always
